# appflowy_kanban/growth-os/deploy/linux/systemd/growthos-curate.service
[Unit]
Description=Growth OS curator (fetch -> score -> dedupe -> upsert)
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=CHANGEME
WorkingDirectory=/home/CHANGEME/appflowy/growth-os
Environment=PYTHONPATH=.
ExecStart=/home/CHANGEME/appflowy/growth-os/.venv/bin/python -m growthos.curate
