# appflowy_kanban/growth-os/deploy/linux/systemd/growthos-brief.timer
[Unit]
Description=Run the Growth OS daily brief every morning

[Timer]
OnCalendar=*-*-* 06:30:00
Persistent=true

[Install]
WantedBy=timers.target
