# archive/appflowy/legacy-growth-os/deploy/linux/systemd/appflowy-backup.service
[Unit]
Description=AppFlowy Postgres + Growth OS state backup
After=docker.service

[Service]
Type=oneshot
User=CHANGEME
ExecStart=/home/CHANGEME/appflowy/growth-os/deploy/linux/backup.sh
